# roles/bbb/templates/microsd.mount.j2
[Unit]
Description=MicroSD flash card

[Mount]
What=/dev/mmcblk0p1
Where=/media/microsd
Type=ext4
Options=rw,nofail,noatime
SloppyOptions=true
DirectoryMode=0755
TimeoutSec=15

