# roles/bbb/templates/microsd.automount.j2
[Unit]
Description=MicroSD flash card automount

[Automount]
Where=/media/microsd
TimeoutIdleSec=301

[Install]
WantedBy=multi-user.target
